# deploy/member-api.service
# member-api systemd service
[Unit]
Description=Meimei Member API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/meimei-member/server
ExecStart=/opt/meimei-member/server/member-api
EnvironmentFile=/opt/meimei-member/server/.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
